# home/dot_config/systemd/user/dropbox.service.d/override.conf
[Service]
# Run the real Dropbox daemon, not /usr/bin/dropbox: the package wrapper can
# spawn /newerversion children and trigger a systemd restart loop.
ExecStart=
ExecStart=/home/buble/.dropbox-dist/dropboxd
KillMode=control-group
